# root/etc/services.d/rclone/finish
#!/command/execlineb -S0
# shellcheck shell=bash

echo "~~~~~~ RCLONE EXITED ~~~~~~"

/run/s6/basedir/bin/halt
